# mate/.config/autostart/xss-lock.desktop
[Desktop Entry]
Type=Application
Name=xss-lock
Exec=xss-lock -- mate-screensaver-command --lock
X-MATE-Autostart-enabled=true
